# _source-projects/clj-vagrant-template/sample/ansible/files/nrepl-upstart.conf
start on started

script
  exec su - vagrant -c 'cd /var/projects/vagrant-sample/;   lein repl :headless :host 0.0.0.0 :port 4242'
end script
